# rootfs/etc/systemd/system/els81-boot-on.service
[Unit]
Description=ELS81 Init Service

[Service]
ExecStart=/usr/bin/env python3 /usr/local/bin/els81-boot-on.py

[Install]
WantedBy=basic.target
