# rootfs/etc/systemd/system/els81-wwan-manager.service
[Unit]
Description=ELS81 Init Service

[Service]
ExecStart=/usr/local/bin/els81-wwan-manager

[Install]
WantedBy=multi-user.target
