# rootfs/etc/systemd/system/multi-func-btn.service
[Unit]
Description=Shutdown Button Service

[Service]
ExecStart=/usr/bin/env python3 /usr/local/bin/funcbtn.py

[Install]
WantedBy=multi-user.target
