# rootfs/etc/systemd/system/zram.service
[Unit]
Description=Zram Init Service

[Service]
ExecStart=/usr/local/bin/zram.sh

[Install]
WantedBy=multi-user.target
